# example/nice-things.crontab
MAILTO=
* * * * * /usr/local/bin/nice-things
